# router-monitor.conf
[program:router-monitor]
command=/usr/bin/python /home/pi/source/python/router-monitor/router-monitor.py
redirect_stderr=true
stdout_logfile=/var/log/supervisor/router-monitor.log
