# router-monitor.conf
[program:router-monitor]
command=/home/pi/.pyenv/shims/python /home/pi/github/router-monitor/router-monitor.py
redirect_stderr=true
stdout_logfile=/var/log/supervisor/router-monitor.log
